# telegram-gpt4o-bot.service
[Unit]
Description=Telegram GPT-4o Bot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/youruser/telegram_gpt4o_bot/bot.py
WorkingDirectory=/home/youruser/telegram_gpt4o_bot
EnvironmentFile=/home/youruser/telegram_gpt4o_bot/.env
Restart=on-failure

[Install]
WantedBy=multi-user.target
